# ecommerce/ansible/roles/magento/templates/media-watch.service
[Unit]
Description=Syncs media files with Rackspace CDN

[Service]
ExecStart=/bin/bash /srv/media.sh
Restart=always
